# etc/systemd/system/bluetooth_battery.timer
[Install]
WantedBy=timers.target

[Unit]
Description=Run bluetooth_battery.service every 1 min

[Timer]
OnCalendar=*:0/1
